# infra/mng-proxy/mng-proxy.service
[Unit]
Description=MNG Kargo API sabit-IP relay
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mng-proxy
Environment=PORT=8080
Environment=PROXY_SECRET=__CHANGE_ME__
Environment=TARGET_BASE_URL=https://api.mngkargo.com.tr
ExecStart=/usr/bin/node /opt/mng-proxy/server.js
Restart=always
RestartSec=3
User=www-data
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
